# python/kodiclient/kodiremote.service
[Unit]
Description=KissMe Kodi remote

[Service]
ExecStart=/usr/local/bin/kodiremote
